# infra/sync-prices.service
[Unit]
Description=Sync market prices (bets page)

[Service]
Type=oneshot
User=nam
WorkingDirectory=/home/nam/nam-website-deploy
EnvironmentFile=/home/nam/nam-website-deploy/.env
ExecStart=/home/nam/.local/bin/uv run python manage.py sync_prices
